# installs/linux/appweb.upstart
#
#   appweb - Embedded web server
#

description	"Embedthis Appweb"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /usr/lib/appweb/bin/appweb --home /etc/appweb

