# package/linux/appweb.openwrt
#!/bin/sh /etc/rc.common
#
#   Startup script for Embedthis Appweb
#   Copyright (c) Embedthis LLC, 2003-2013. All Rights Reserved.
#
START=50
STOP=50
APPWEB_CMD=${prefixes.bin}/${settings.product}
APPWEB_DIR=${prefixes.cfg}

start() {
    start-stop-daemon -S -b -x "${APPWEB_CMD}" -- --home "$APPWEB_DIR"
}

stop() {
    start-stop-daemon -K -x "${APPWEB_CMD}" >/dev/null 2>&1 ; true
}
